# extra-networks/tools/wifi-recover
#!/bin/sh /etc/rc.common
# Recovers WiFi radios whose VAPs were lost to the mac80211 race condition
# triggered by wifi reload on MT7986. Safe on any OpenWrt machine: if VAPs
# are already up the script exits without doing anything.
START=99

start() {
    sleep 3
    for _hconf in /var/run/hostapd-*.conf; do
        [ -f "$_hconf" ] || continue
        _phy="${_hconf##*/hostapd-}"; _phy="${_phy%.conf}"
        iw dev 2>/dev/null | grep -q "Interface ${_phy}-" && continue
        logger -t wifi-recover "${_phy}: no VAPs up, recovering via hostapd config_set"
        ubus call hostapd config_set \
            "{\"phy\":\"${_phy}\",\"radio\":-1,\"config\":\"${_hconf}\",\"prev_config\":\"${_hconf}.prev\"}" \
            >/dev/null 2>&1 || true
    done
}
